# misc/procs/com.deepin.system.procs.service
[Service]
Name=com.deepin.system.procs
Exec=/usr/lib/deepin-daemon/netlink
