# misc/procs/com.deepin.system.procs.service
[Service]
Name=com.deepin.system.Procs
Exec=/home/aris/Desktop/WorkProject/DeepinProxy/Bin/NetLink/NetLink
